# services/pp-llm.service
[Unit]
Description=PinePhone Voice Agent — LLM server (llama.cpp)
After=default.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
Environment=LD_LIBRARY_PATH=%h/llama.cpp/build/bin
ExecStart=%h/llama.cpp/build/bin/llama-server \
    -m %h/models/gemma-3-270m-it-Q4_K_M.gguf \
    -t 4 --ctx-size 512 --reasoning-budget 0 \
    --flash-attn on --mlock --prio 2 \
    --log-disable
Restart=on-failure
RestartSec=3
AllowedCPUs=4-5
NoNewPrivileges=yes

[Install]
WantedBy=default.target
